# supervisord.conf
[supervisord]
nodaemon=true

[program:logstash]
autostart=true
autorestart=true
startsecs=10
command=/etc/init.d/logstash start
